# deploy/systemd/trademaster-web.service
[Unit]
Description=TradeMaster Django (Gunicorn)
After=network.target redis-server.service
Wants=redis-server.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/trademaster/backend
Environment="PATH=/var/www/trademaster/backend/venv/bin"
ExecStart=/var/www/trademaster/backend/venv/bin/gunicorn \
  --workers 2 \
  --bind 127.0.0.1:8000 \
  --timeout 120 \
  trademaster_project.wsgi:application
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
